# deploy/supervisor/heyvmd.conf
; supervisord program config for heyvmd, the heyo VM control-plane daemon
; pg-vm-pool depends on (serves the local sandbox HTTP API on 34099 — see
; pg-vm-pool.conf, which talks to it over plain HTTP at 127.0.0.1:34099 and
; needs it running first).
;
; Install by copying (or symlinking) into your supervisor include dir, e.g.:
;   sudo mkdir -p /var/log/heyvmd && sudo chown sam:sam /var/log/heyvmd
;   sudo ln -s /home/sam/Projects/pg-fc/deploy/supervisor/heyvmd.conf \
;              /etc/supervisor/conf.d/heyvmd.conf
;   sudo supervisorctl reread && sudo supervisorctl update
; then:  supervisorctl {start,stop,restart,tail} heyvmd
;
; This installs TWO programs: heyvmd itself, and a health-check sidecar that
; force-restarts it if its API stops responding (see heyvmd-healthcheck.sh).

[program:heyvmd]
command=/usr/bin/heyvmd --api-port 34099
directory=/home/sam
user=sam
; supervisord does NOT inherit the login shell's PATH — children get a minimal
; environment. heyvm invokes `firecracker` by bare name (relying on PATH), and
; firecracker lives in /usr/local/bin, so without this it fails with
; "Firecracker is not available". Include /usr/local/sbin + /usr/sbin too since
; VM setup shells out to sbin tools (ip, useradd, etc.).
environment=PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",HOME="/home/sam",USER="sam"
autostart=true
autorestart=true
startsecs=5
startretries=3
stopsignal=TERM
stopwaitsecs=15
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/heyvmd/heyvmd.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5

; --- health check / watchdog ---
; `autorestart` above only fires once the heyvmd *process* exits; it does
; nothing if the process is alive but wedged (API stopped responding while
; the daemon itself is still up). This sidecar polls GET /deployed-sandboxes
; every HEYVMD_HEALTHCHECK_INTERVAL_SECS and, after FAIL_THRESHOLD
; consecutive failures, runs `supervisorctl restart heyvmd`.
;
; Deliberately has no `user=` (so it runs as whatever supervisord itself runs
; as, normally root): /var/run/supervisor.sock is 0700 root-owned, so a
; user=sam process cannot call supervisorctl at all.
[program:heyvmd-healthcheck]
command=/home/sam/Projects/pg-fc/deploy/supervisor/heyvmd-healthcheck.sh
autostart=true
autorestart=true
startsecs=5
environment=
    HEYVMD_HEALTHCHECK_PORT="34099",
    HEYVMD_HEALTHCHECK_INTERVAL_SECS="15",
    HEYVMD_HEALTHCHECK_FAIL_THRESHOLD="3",
    HEYVMD_HEALTHCHECK_COOLDOWN_SECS="60"
redirect_stderr=true
stdout_logfile=/var/log/heyvmd/heyvmd-healthcheck.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
