# deploy/supervisor/pg-vm-pool.conf
; supervisord program config for the pg-vm-pool Postgres pooler.
;
; Install by copying (or symlinking) into your supervisor include dir, e.g.:
;   sudo ln -s /home/sam/Projects/pg-fc/deploy/supervisor/pg-vm-pool.conf \
;              /etc/supervisor/conf.d/pg-vm-pool.conf
;   sudo supervisorctl reread && sudo supervisorctl update
; then:  supervisorctl {start,stop,restart,tail} pg-vm-pool
;
; Manage all runtime config from the `environment=` line below — edit it, then
; `supervisorctl reread && supervisorctl update pg-vm-pool` (a plain restart
; does NOT reload environment=; update does).
;
; Build the release binary first:  cargo build --release
; The pooler needs a local heyvmd (POST /sandboxes/:id/tcp-tunnel) already
; running; run that under its own supervisor program if it isn't a system svc.

[program:pg-vm-pool]
command=/home/sam/Projects/pg-fc/target/release/pg-vm-pool
directory=/home/sam/Projects/pg-fc
user=sam

; --- runtime config (this is the point of the unit: manage env here) ---
; Values with commas (e.g. KEEPALIVE_SCHEMAS) MUST be double-quoted so the
; comma isn't read as a field separator. Omit a var to take the code default.
;   PG_VM_POOL_LISTEN                default 127.0.0.1:6432
;   PG_VM_POOL_IMAGE                 default pg
;   PG_VM_POOL_USER                  default postgres
;   PG_VM_POOL_PASSWORD              default unset (no client auth gate; trust
;                                    backend auth). Set to require this exact
;                                    password from clients (cleartext challenge,
;                                    rejects anyone who doesn't send it) AND to
;                                    use it for the probe/bootstrap if the
;                                    backend needs scram/md5 instead of trust.
;                                    Pair with TLS below unless LISTEN is loopback.
;   PG_VM_POOL_IDLE_TIMEOUT_SECS     default 900; 0 disables idle reaping
;   PG_VM_POOL_READY_TIMEOUT_SECS    default 300
;   PG_VM_POOL_CONNECT_TIMEOUT_SECS  default 30 (iroh tunnel handshake cap)
;   PG_VM_POOL_DIRECT_CONNECT        default on; "0" forces the iroh tunnel path
;   PG_VM_POOL_DATA_DISK_GB          default 4
;   PG_VM_POOL_KEEPALIVE_SCHEMAS     default none; comma-separated, quote it
;   PG_VM_POOL_STATE_FILE            default ~/.heyo/pg-vm-pool/registry.tsv
;   PG_VM_POOL_TLS_CERT              default unset (TLS off); PEM cert chain path
;   PG_VM_POOL_TLS_KEY               default unset; PEM private key path — set
;                                    BOTH or NEITHER; files hot-reload on change
;                                    (certbot renewals need no restart)
;   RUST_LOG                         default info,pg_vm_pool=info
environment=
    PG_VM_POOL_LISTEN="127.0.0.1:6432",
    PG_VM_POOL_IMAGE="pg",
    PG_VM_POOL_USER="postgres",
    PG_VM_POOL_IDLE_TIMEOUT_SECS="300",
    PG_VM_POOL_READY_TIMEOUT_SECS="300",
    PG_VM_POOL_DATA_DISK_GB="4",
    PG_VM_POOL_KEEPALIVE_SCHEMAS="",
    RUST_LOG="info,pg_vm_pool=info",
    PG_VM_POOL_PASSWORD="secret"
; To enable client-facing TLS, append (see README "TLS" section):
;    PG_VM_POOL_TLS_CERT="/home/sam/.heyo/pg-vm-pool/tls/fullchain.pem",
;    PG_VM_POOL_TLS_KEY="/home/sam/.heyo/pg-vm-pool/tls/privkey.pem"

; --- process management ---
autostart=true
autorestart=true
startsecs=5
startretries=3
stopsignal=TERM
stopwaitsecs=15
stopasgroup=true
killasgroup=true

; --- logs ---
redirect_stderr=true
stdout_logfile=/var/log/pg-vm-pool/pg-vm-pool.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
